# MY_TEST/scripts/com.mytest.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.mytest.backup</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-lc</string>
      <string>cd /Users/ivangritsenko/MY_TEST/MY_TEST && ./scripts/backup-data.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>0</integer>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/ivangritsenko/MY_TEST/MY_TEST/backups/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ivangritsenko/MY_TEST/MY_TEST/backups/launchd.error.log</string>
  </dict>
</plist>
